# conf/emailpoubelle.cron
0 */2 * * * __APP__ cd __FINALPATH__/www/; /usr/bin/php index.php > /dev/null 2>&1