# conf/emailpoubelle.cron
0 */2 * * * www-data cd $install_dir/www/; /usr/bin/php index.php > /dev/null 2>&1